# config/etc/rc.local
#!/bin/sh
set -x

chmod u=rwx,go= /root
chmod u=rwx,go= /root/.ssh
chmod u=rw,go= /root/.ssh/authorized_keys
echo "hardenedbsd" | pw usermod root -h 0
